# etc/supervisor/geigor.conf
[program:geigor]
logfile=/var/log/geigor.log
command=/path/to/geigor/geigor http://your-domain.com:9319/ 9319 index demo

autostart=true
autorestart=true
startsecs=10
startretries=3